# .config/systemd/user/devilspie.service
[Unit]
Description=Devilspie

[Service]
ExecStart=/usr/bin/devilspie

[Install]
WantedBy=default.target
