# .config/systemd/user/irssi-libnotify.service
[Unit]
Description=irssi-libnotify

[Service]
ExecStart=python2.7 /usr/bin/notify-listener

[Install]
WantedBy=default.target
