# .config/systemd/user/irssi-libnotify.service
[Unit]
Description=irssi-libnotify-git

[Service]
ExecStart=/usr/bin/notify-listener

[Install]
WantedBy=default.target
