# .config/systemd/user/pia-vpn.service
[Unit]
Description=PIA VPN

[Service]
ExecStart=/home/n0t/scripts/vpn/./nm-pia.sh "PIA - CA Montreal"

[Install]
WantedBy=default.target

